# deployment/cloud/dashboard/urac_dashboard/FILES/supervisord.conf
[supervisord]
nodaemon=true
environment=SOAJS_ENV="prod",SOAJS_PRJ="default/",SOAJS_REGDIR="/opt/soajs/FILES/"

[program:urac]
command=/usr/bin/node /opt/soajs/node_modules/soajs.urac/index.js
autostart=true
autorestart=true
stderr_logfile=/opt/soajs/log/soajs.urac.err.log
stdout_logfile=/opt/soajs/log/soajs.urac.out.log

[program:dashboard]
command=/usr/bin/node /opt/soajs/node_modules/soajs.dashboard/service/index.js
autostart=true
autorestart=true
stderr_logfile=/opt/soajs/log/soajs.dashboard.err.log
stdout_logfile=/opt/soajs/log/soajs.dashboard.out.log