# deployment/cloud/dashboard/controller/FILES/supervisord.conf
[supervisord]
nodaemon=true
environment=SOAJS_ENV="prod",SOAJS_PRJ="default/",SOAJS_REGDIR="/opt/soajs/FILES/"

[program:controller]
command=/usr/bin/node /opt/soajs/node_modules/soajs.controller/index.js
autostart=true
autorestart=true
stderr_logfile=/opt/soajs/log/soajs.controller.err.log
stdout_logfile=/opt/soajs/log/soajs.controller.out.log
